# deploy/langflow-chat.service
[Unit]
Description=AlfyAI Application
After=network.target

[Service]
Type=simple
User=alfyai
Group=alfyai
WorkingDirectory=/opt/alfyai
EnvironmentFile=/opt/alfyai/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/alfyai/build/index.js
Restart=on-failure
RestartSec=10
TimeoutStopSec=90
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
